# rabbitmq/supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]
port=9001

[supervisord]
logfile=/var/log/supervisord.log ; (main log file)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
user=root                    ; (default is current user, required if root)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
autostart=true
user=root
startsecs=3                   ; number of secs prog must stay running (def. 1)
redirect_stderr=true
stdout_logfile=/var/log/rabbitmq.log

[program:ssh]
command=/usr/sbin/sshd -D
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
autostart=true
user=root
startsecs=3                   ; number of secs prog must stay running (def. 1)
redirect_stderr=true
stdout_logfile=/var/log/sshd.log
